# files/rc.local
echo XHC > /proc/acpi/wakeup
echo GLAN > /proc/acpi/wakeup
